# templates/kopia-backup.timer.j2
# {{ ansible_managed }}

# SPDX-License-Identifier: MIT
[Unit]
Description=Kopia backup timer

[Timer]
OnCalendar={{ kopia_backup_execution_timer }}
Persistent=true

[Install]
WantedBy=timers.target