# templates/kopia-backup.service.j2
# SPDX-License-Identifier: MIT
[Unit]
Description=Kopia backup service

[Service]
Type=oneshot
ExecStart="{{ kopia_scripts_dir }}/kopia-job.sh"